# rootfs/etc/s6-overlay/s6-rc.d/iredapd/run
#!/command/execlineb -P
with-contenv

foreground { echo "Starting iRedAPD..." }

# Start iRedAPD
cd /opt/iredapd
python3 /opt/iredapd/iredapd.py
