# rootfs/etc/s6-overlay/s6-rc.d/iredapd/run
#!/command/with-contenv bash

echo "Starting iRedAPD..."

# Create runtime directory
mkdir -p /run/iredapd /var/log/iredapd
chown iredapd:iredapd /run/iredapd /var/log/iredapd 2>/dev/null || true

# Clean stale PID file
rm -f /run/iredapd/iredapd.pid

# Small delay to allow socket cleanup on restart
sleep 1

# Start iRedAPD in foreground
cd /opt/iredapd
exec python3 /opt/iredapd/iredapd.py --foreground
