# rootfs/etc/s6-overlay/s6-rc.d/rsyslog/run
#!/command/execlineb -P
with-contenv
/usr/sbin/rsyslogd -n
