# rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv

foreground { echo "Starting PHP-FPM..." }

# Ensure run directory exists
foreground { mkdir -p /run/php }
foreground { chown www-data:www-data /run/php }

# Start php-fpm in foreground
/usr/sbin/php-fpm8.1 -F
