# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv

foreground { echo "Starting Nginx..." }

# Test configuration
foreground { nginx -t }

# Start nginx in foreground
nginx -g "daemon off;"
