# rootfs/etc/s6-overlay/s6-rc.d/cert-reload/run
#!/command/execlineb -P
with-contenv

foreground { echo "Starting certificate reload watcher..." }

/etc/s6-overlay/scripts/cert-reload.sh
