# rootfs/etc/s6-overlay/s6-rc.d/clamav/run
#!/command/with-contenv bash
set -e

# Check if ClamAV is disabled
if [ "${DISABLE_CLAMAV}" = "YES" ]; then
    echo "ClamAV disabled, sleeping..."
    exec sleep infinity
fi

echo "Starting ClamAV daemon..."

# Create directories and ensure proper ownership
mkdir -p /run/clamav /var/log/clamav /var/lib/clamav
chown -R clamav:clamav /run/clamav /var/lib/clamav /var/log/clamav

# Clean stale files
rm -f /run/clamav/clamd.pid /run/clamav/clamd.ctl

# Wait for any other freshclam process to finish
while pgrep -x freshclam > /dev/null 2>&1; do
    echo "Waiting for freshclam to finish..."
    sleep 2
done

# Download virus definitions if not present
if [ ! -f "/var/lib/clamav/main.cvd" ]; then
    echo "Virus definitions not found, downloading..."
    su -s /bin/bash -c "freshclam --foreground --stdout" clamav || true
    echo "Virus definitions downloaded, starting clamd..."
fi

echo "Starting clamd..."
exec /usr/sbin/clamd --foreground=true
